# deploy/devmail/proidentity-webadmin.service
[Unit]
Description=ProIdentity Mail Webadmin
After=network-online.target mariadb.service
Wants=network-online.target

[Service]
Type=simple
User=proidentity
Group=proidentity
EnvironmentFile=/etc/proidentity-mail/proidentity-mail.env
ExecStart=/opt/proidentity-mail/bin/webadmin
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/etc/proidentity-mail

[Install]
WantedBy=multi-user.target
