# deploy/devmail/proidentity-backup.timer
[Unit]
Description=Run ProIdentity Mail Backup Daily

[Timer]
OnCalendar=*-*-* 02:15:00
RandomizedDelaySec=20m
Persistent=true
Unit=proidentity-backup.service

[Install]
WantedBy=timers.target
